# systemd/pilot-room-agent.service.j2
[Unit]
Description=Pilot room endpoint agent
After=network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=simple
User={{ room_endpoint_user }}
Group={{ room_endpoint_user }}
Environment=XDG_RUNTIME_DIR=/run/user/{{ room_endpoint_uid.stdout }}
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/{{ room_endpoint_uid.stdout }}/bus
{% if room_endpoint_video_enabled | bool and room_endpoint_video_wayland_display | length > 0 %}
Environment=WAYLAND_DISPLAY={{ room_endpoint_video_wayland_display }}
{% endif %}
Environment=PILOT_CONFIG=/etc/pilot/room.toml
Environment=PYTHONPATH=/opt/pilot/current/apps/room-agent
ExecStart=/opt/pilot/current/.venv/bin/python -m pilot_room_agent.main
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/pilot {{ room_endpoint_audio_cache_path }}

[Install]
WantedBy=multi-user.target
